# aprs-uptime.service
[Unit]
Description=APRS Nearby Watch & Kuma Push
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/APRS_Uptime
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/APRS_Uptime/.venv/bin/python /home/ubuntu/APRS_Uptime/aprs_nearby_watch.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
